# admin1/task3/bg_script.service
[Unit]
Description=Lazy hello script daemon
[Service]
ExecStart=/bin/bash /usr/local/bin/script1.sh
[Install]
WantedBy=multi-user.target

